# root/app/dnsmasq/supervisord.template.conf
[program:dnsmasq]
command = /usr/sbin/dnsmasq -d -k -q -C /app/dnsmasq/dnsmasq.conf
autorestart = true
priority = 1
startsecs = 1
stdout_logfile=/var/log/dnsmasq.log
logfile_maxbytes = 50MB
