# root/app/smartdns/supervisord.conf
[program:sniproxy]
command = /usr/sbin/sniproxy -f -c /etc/sniproxy/sniproxy.conf
autorestart = true
pidfile = /run/sniproxy.pid
priority = 2
